# conf/scrape.cron
* * * * * /opt/scrapy/scrape.sh >> /var/log/cron.log 2>&1
